# profile.d/wdd130.sh
#!/usr/bin/bash.exe
echo "Load wdd130 profile!"